# deploy/alou-pay.service
[Unit]
Description=Alou Pay Rust Backend Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/alou-pay/backend
ExecStart=/opt/alou-pay/backend/agent-http-server
Restart=always
RestartSec=10
Environment=PORT=3001
Environment=RUST_LOG=info
EnvironmentFile=-/opt/alou-pay/backend/.env

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/alou-pay /var/log/alou-pay

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alou-pay

[Install]
WantedBy=multi-user.target
